# jupiter-current/lib/udev/rules.d/80-drm.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/lib/jupiter/scripts/vga-out mon"
